# ai-codeagent.service
[Unit]
Description=AI Code Agent (FastAPI + uvicorn)
After=network.target

[Service]
Type=simple
User=user
WorkingDirectory=/home/user/AI-Codeagent
ExecStart=/home/user/AI-Codeagent/venv/bin/uvicorn server:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
EnvironmentFile=/home/user/AI-Codeagent/.env

[Install]
WantedBy=multi-user.target
